# homebin/rc.local
#!/bin/bash

/sbin/modprobe vboxdrv

exec setuidgid yuting /usr/bin/VBoxManage startvm winxp --type headless

