# homebin/rc.local
#!/bin/bash

/sbin/modprobe vboxdrv

exec /usr/bin/VBoxManage startvm winxp --type headless

