# homebin/rc.local
#!/bin/bash

/sbin/modprobe vboxdrv

cd /home/yuting
exec setuidgid yuting /usr/bin/VBoxManage startvm winxp --type headless

