# homebin/rc.local
#!/bin/bash

/sbin/modprobe vboxdrv

sudo -u yuting sh -c "/usr/bin/VBoxManage startvm winxp --type headless"

